# sv/ssh-agent/run
#!/bin/zsh

SSH_AGENT_ENV=$HOME/.ssh/ssh-agent-env
exec > "$SSH_AGENT_ENV"
exec ssh-agent -D